# 7.2/image-files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm -R
priority=999

[program:nginx]
command=nginx -g "daemon off;"
priority=10
